# crontab/dedeco99
0,30 * * * * polybar-msg hook updates 1
