# crontab/dedeco99
0,30 * * * * polybar-msg action "#updates.hook.0"
